# config/tkmcollector.service.in
###############################################################################
#  Task Monitor Collector version @PROJECT_VERSION@
###############################################################################
[Unit]
Description=TaskMonitor Collector Service

[Service]
Type=simple
User=root
Group=root
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/tkmcollector -c /etc/tkmcollector.conf
WatchdogSec=30

[Install]
WantedBy=multi-user.target
